# crontab
# 设置环境变量
COOKIE="98堂ck"
TG_USER_ID="your_tg_user_id"
TG_BOT_TOKEN="your_tg_bot_token"
PROXY_URL="http://your.proxy.server:port"

# 每天定时执行 daysign.py 脚本
0 0 * * * python /app/daysign.py
